# templates/systemd.service.tmpl
[Unit]
Description=A privacy respecting, hackable metasearch engine

[Service]
ExecStart=/usr/bin/searx-run

[Install]
WantedBy=default.target
